# supervisord-logstash.conf
[program:logstash-agent]
command=/start-logstash-agent.sh
autostart=true
auto_restart=true

[program:logstash-web]
command=/start-logstash-web.sh
autostart=true
autorestart=true
