# supervisord-elasticsearch.conf
[program:elasticsearch]
command=/start-elasticsearch.sh
autostart=true
auto_restart=true
